# rootfs/etc/crontabs/root
*/15 * * * * sleep $(($RANDOM % (5 * 60))); /etc/update-mydns
